# deploy/systemd/sunmessenger-pip-audit.timer
[Unit]
Description=Run SUN Messenger pip-audit weekly

[Timer]
OnCalendar=weekly
Persistent=true
Unit=sunmessenger-pip-audit.service

[Install]
WantedBy=timers.target

